# templates/apt-daily-upgrade.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ untnd_pgrds__apt_daily_upgrade_timer_on_calendar }}
RandomizedDelaySec={{ untnd_pgrds__apt_daily_upgrade_timer_randomized_delay }}
