# templates/apt-daily.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ untnd_pgrds__apt_daily_timer_on_calendar }}
RandomizedDelaySec={{ untnd_pgrds__apt_daily_timer_randomized_delay }}
